# daemon/launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this service -->
    <key>Label</key>
    <string>com.claude.telegram-relay</string>

    <!-- Command to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME/.claude-relay/start.sh</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/dev/claude-telegram-relay</string>

    <!-- Start when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep running (restart if it crashes) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log files -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/claude-telegram-relay.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/claude-telegram-relay.error.log</string>

    <!-- Wait 10 seconds before restarting after crash -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>

<!--
SETUP INSTRUCTIONS:

1. Copy to LaunchAgents:
   cp daemon/launchagent.plist ~/Library/LaunchAgents/com.claude.telegram-relay.plist

2. Load the agent:
   launchctl load ~/Library/LaunchAgents/com.claude.telegram-relay.plist

3. Check if running:
   launchctl list | grep claude

4. View logs:
   tail -f ~/Library/Logs/claude-telegram-relay.log

5. To stop:
   launchctl unload ~/Library/LaunchAgents/com.claude.telegram-relay.plist

6. To restart:
   launchctl unload ~/Library/LaunchAgents/com.claude.telegram-relay.plist
   launchctl load ~/Library/LaunchAgents/com.claude.telegram-relay.plist
-->
